# services/greenhouse-ai-decision.service
[Unit]
Description=Greenhouse Auto IA — Claude decide actuadores (oneshot, triggered by timer)
After=network-online.target greenhouse.service
Wants=network-online.target

[Service]
Type=oneshot
User=kowen
Group=kowen
WorkingDirectory=/home/kowen/proyectos/Control_clima
EnvironmentFile=/home/kowen/proyectos/Control_clima/.env
ExecStart=/home/kowen/proyectos/Control_clima/.venv/bin/python /home/kowen/proyectos/Control_clima/bin/ai-decide.py
TimeoutStartSec=120
StandardOutput=append:/var/log/greenhouse-ai-decision.log
StandardError=append:/var/log/greenhouse-ai-decision.log
